# services/incbak-root.service
[Unit]
Description=Incremantal backup for root

[Service]
Type=oneshot
ExecStart=/mnt/btrfs/incbak-root
