# services/incbak-borg.service
[Unit]
Description=Incremantal backup for borg

[Service]
User=mitch
Type=oneshot
ExecStart=/mnt/btrfs/incbak-borg
