# services/incbak-borg.timer
[Unit]
Description=Incremantal backup for borg

[Timer]
OnCalendar=*-*-* 03:30:00

[Install]
WantedBy=timers.target

