# services/incbak-zfs.service
[Unit]
Description=Incremantal backup for zfs

[Service]
Type=oneshot
ExecStart=/mnt/btrfs/incbak-zfs
