# services/warmcache.service
[Unit]
Description=Warm disk cache

[Service]
IOSchedulingClass=idle
Type=oneshot
ExecStart=/home/mitch/.bin/warmcache
