# services/warmcache.timer
[Unit]
Description=Warm disk cache

[Timer]
OnBootSec=120

[Install]
WantedBy=timers.target

