# services/warmcache.service
[Unit]
Description=Warm disk cache
ConditionVirtualization=!vm

[Service]
IOSchedulingClass=idle
Type=oneshot
ExecStart=/home/mitch/.bin/warmcache
